# data/systemd/com.system76.CosmicPaste.service
[Unit]
Description=COSMIC Paste clipboard daemon
Documentation=https://github.com/erik-balfe/cosmic-paste
After=graphical-session.target

[Service]
Type=dbus
BusName=org.system76.CosmicPaste
ExecStart=@bindir@/cosmic-paste-daemon
Restart=on-failure
TimeoutStopSec=10
KillMode=mixed
# Inherit WAYLAND_DISPLAY / XDG_RUNTIME_DIR from the graphical session.
ImportEnvironment=WAYLAND_DISPLAY XDG_RUNTIME_DIR DISPLAY

[Install]
WantedBy=graphical-session.target